# deployment/systemd/voiceforpalestine-backend.service
[Unit]
Description=Voice for Palestine Backend Service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/voiceforpalestine.xyz/backend
ExecStart=/usr/bin/node server.js
Restart=on-failure
Environment=NODE_ENV=production
Environment=PORT=3001

[Install]
WantedBy=multi-user.target
